# selfhosted-docs.socket
[Unit]
Description=Selfhosted Docs Socket
PartOf=selfhosted-docs.service

[Socket]
ListenStream=127.0.0.1:9997

[Install]
WantedBy=sockets.target
